# docker/crontab
* * * * *   /opt/app/gen-releases-task.sh
